# root/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P
with-contenv
s6-setuidgid root
nginx -c /config/nginx/nginx.conf

